# svisitor.ini
[program:visitor_web]
command=/var/www/html/visitor.videonow.ru/golang/visitor/visitor-linux-x64 -serve=web -cpu=4
numprocs=1
directory=/var/www/html/visitor.videonow.ru/
stdout_logfile=/var/log/visitor_web.log
autostart=true
autorestart=true
user=apache
stopsignal=KILL